# tech-radar.conf
description "tech-radar"
author "abtv <github.com/abtv>"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Essentially lets upstart know the process will detach itself to the background
# This option does not seem to be of great importance, so it does not need to be set.
#expect fork

# Specify working directory
chdir /root/radar/

# Specify the process/command to start, e.g.
exec java -jar backend.jar
